# supervisor/elasticsearch.ini
[program:elasticsearch]
command=/opt/elasticsearch/bin/elasticsearch --path.home=/opt/elasticsearch/ --path.data=/var/elasticsearch/data/ --path.logs=/var/elasticsearch/logs 
stdout_logfile=/var/log/elasticsearch.stdout.log
stderr_logfile=/var/log/elasticsearch.stderr.log

