# supervisor/logstash.ini
[program:logstash]
command=/opt/logstash/bin/logstash -f /etc/logstash.d
stdout_logfile=/var/log/logstash.stdout.log
stderr_logfile=/var/log/logstash.stderr.log

